# rtc_alarm_scripts/rm_wakeup.timer
[Unit]
Description=Re-arm the RTC alarm for 04:00 Central after every boot or wake

[Timer]
OnBootSec=30s
OnWakeSec=30s      # Any reasonable period ≥ the wake interval is fine
Unit=rm_wakeup.service

[Install]
WantedBy=multi-user.target
